# supervisord.conf
[program:php]
command=/usr/sbin/php-fpm7.2-zts

[program:filebeat]
command=/app/filebeat.sh